# etc/systemd/system/wifiacc.service
[Unit]
Description=Makes wifi usable on HP

[Service]
ExecStart=/usr/bin/wifiacc

[Install]
WantedBy=multi-user.target
